# workflow/pg_migration/vm_setup/cloud-sql-proxy.service
[Unit]
Description=Cloud SQL Auth Proxy (mes-pg-staging)
Documentation=https://cloud.google.com/sql/docs/postgres/sql-proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cloudsqlproxy
Group=cloudsqlproxy
WorkingDirectory=/home/cloudsqlproxy
# The Go-based proxy reads several gcloud SDK config files under
# /home/$USER/.config/gcloud and /home/$USER/.secureConnect via Go
# stdlib os/user (which resolves via getpwuid + a hardcoded /home/$USER
# fallback) — Environment=HOME does NOT redirect those reads. Easiest
# fix is to provision the home dir for real (install.sh does this) and
# leave ProtectHome off so the proxy can reach it.
# --private-ip forces the proxy to use the VPC private IP path we set up
# with VPC peering, NOT the public internet. The proxy still mediates IAM
# auth + encryption — we get connection-string-style ergonomics
# (psql -h 127.0.0.1 -p 5432) without exposing Postgres on a public IP.
ExecStart=/opt/cloud-sql-proxy/cloud-sql-proxy \
    --address=127.0.0.1 \
    --port=5432 \
    --private-ip \
    superb-metric-492315-r5:us-central1:mes-pg-staging
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening. ProtectHome dropped — see comment above.
# ProtectSystem=full keeps /usr /boot /etc read-only but allows /var, /tmp,
# and /home (which is what we need).
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
